# composer/siem/sophos/config/crontab
* * * * * cd /app/sophos && python3 siem.py >> /var/log/sophos/sophos-$(date +'%Y-%m-%d').log
