# __template/config/supervisord.conf
[supervisord]
pidfile=/var/run/supervisord.pid
logfile=/dev/null
logfile_backups=0
logfile_maxbytes=0
loglevel=info
nodaemon=true
user=root

[program:docker-entrypoint]
command=/bin/bash -c "/usr/local/bin/docker-entrypoint.sh"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
